# dotfiles-extras/etc/local.d/keymap.start
#!/bin/sh
loadkeys /etc/esc-caps-swap.map
